# mail/Library/LaunchAgents/homebrew.mxcl.offline-imap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <false/>
    <key>Label</key>
    <string>homebrew.mxcl.offline-imap</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/jterk/bin/mailrun.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
      <!-- python2 is in /usr/local/bin -->
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jterk/tmp/mailrun.stdout</string>
    <key>StandardErrorPath</key>
    <string>/Users/jterk/tmp/mailrun.stderr</string>
  </dict>
</plist>
